# common/component/superset/files/superset.service.template
#superset.service
###########
[Unit]
Description=Superset Visualizaton Platform

[Service]
Type=simple


ExecStart=/etc/init.d/superset


[Install]
WantedBy=default.target
